# motion/crontab
* * * * * rm -rf %%TARGETDIR%%/*
* * * * * rm -rf /share/motion/camera0
